# memoryd/src/memoryd/templates/systemd-weekly-identity.timer.j2
[Unit]
Description=memoryd weekly identity rewrite timer

[Timer]
OnCalendar=Sun *-*-* {{ "%02d:%02d:00" | format(hour, minute) }}
Persistent=true

[Install]
WantedBy=timers.target
