# memoryd/src/memoryd/templates/systemd-weekly-identity.service.j2
[Unit]
Description=memoryd weekly identity rewrite
After=default.target

[Service]
Type=oneshot
ExecStart={{ memoryd_bin }} profile rewrite
Environment="MEMORYD_DATA_ROOT={{ data_root }}"
StandardOutput=append:{{ log_dir }}/{{ label }}.out.log
StandardError=append:{{ log_dir }}/{{ label }}.err.log
