# memoryd/src/memoryd/templates/systemd-decay.service.j2
[Unit]
Description=memoryd daily decay sweep
After=default.target

[Service]
Type=oneshot
ExecStart={{ memoryd_bin }} decay-sweep
Environment="MEMORYD_DATA_ROOT={{ data_root }}"
StandardOutput=append:{{ log_dir }}/{{ label }}.out.log
StandardError=append:{{ log_dir }}/{{ label }}.err.log
