# memoryd/src/memoryd/templates/systemd-sync-push.service.j2
[Unit]
Description=memoryd daily sync push (export to sync.dir)
After=default.target

[Service]
Type=oneshot
ExecStart={{ memoryd_bin }} sync export --auto
Environment="MEMORYD_DATA_ROOT={{ data_root }}"
StandardOutput=append:{{ log_dir }}/{{ label }}.out.log
StandardError=append:{{ log_dir }}/{{ label }}.err.log
