# memoryd/src/memoryd/templates/systemd-decay.timer.j2
[Unit]
Description=memoryd daily decay timer

[Timer]
OnCalendar=*-*-* {{ "%02d:%02d:00" | format(hour, minute) }}
Persistent=true

[Install]
WantedBy=timers.target
