# scripts/launchd/com.memoryd.mirror.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.memoryd.mirror</string>
  <key>ProgramArguments</key>
  <array>
    <string>__MEMORYD_BIN__</string>
    <string>mirror</string>
    <string>--codex</string>
    <string>--openclaw</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>MEMORYD_DATA_ROOT</key>
    <string>__MEMORYD_DATA_ROOT__</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>
  <key>StandardOutPath</key>
  <string>__MEMORYD_DATA_ROOT__/logs/mirror.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__MEMORYD_DATA_ROOT__/logs/mirror.stderr.log</string>
</dict>
</plist>
